# deploy/systemd/certwatch-daemon.service
# certwatch as a long-running daemon (the `run` subcommand sleeps `interval`
# between cycles). Prefer this over the timer when you want certwatch to own
# its own schedule, e.g. in a container or on a host without systemd timers.
#
# Install:
#   sudo cp deploy/systemd/certwatch-daemon.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now certwatch-daemon.service

[Unit]
Description=certwatch — continuous TLS certificate monitor
Documentation=https://github.com/mangrisano/certwatch
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=certwatch
Group=certwatch
ExecStart=/usr/local/bin/certwatch run -c /etc/certwatch/certwatch.yml
Restart=on-failure
RestartSec=30

StateDirectory=certwatch
StateDirectoryMode=0750

# Hardening (see certwatch.service for rationale).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
LockPersonality=true

[Install]
WantedBy=multi-user.target
